# run_print_controller.service
[Unit]
Description=Run Docker Compose for AR Print Controller
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
RemainAfterExit=yes
WorkingDirectory=/home/kicks-kiosk/ar-print-controller
ExecStartPre=/usr/local/bin/docker-compose down -v
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down -v
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target

